# api/supervisor.conf
[supervisord]
nodaemon=true

[program:api]
command=/usr/local/bin/api > /proc/1/fd/1
directory=/opt/go/pesto/api/
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
autostart=true
stopsignal=QUIT
stopwaitsecs=120

